# systemd/.config/systemd/user/kbuildsycoca6.service
[Unit]
Description=Rebuild KDE system configuration cache
After=graphical-session.target
PartOf=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=oneshot
ExecStart=/usr/bin/kbuildsycoca6
Slice=session.slice

[Install]
WantedBy=graphical-session.target
